# conf/supervisor/omxremote.conf
[program:omxremote]
command=gunicorn -k eventlet -w 2 -b 0.0.0.0:5000 omxremote:app
autostart=true
autorestart=true
startretries=3
user=pi

